# pkg/debian/postinst
#!/bin/sh
set -e

echo " ┌──────────────────────────────────────────────────────────┐"
echo " │                      datawall                           │"
echo " │  Network data usage tracker per program                 │"
echo " ├──────────────────────────────────────────────────────────┤"
echo " │                                                         │"
echo " │  To enable the daemon on boot:                          │"
echo " │    systemctl --user enable --now datawall.service       │"
echo " │                                                         │"
echo " │  Commands:                                              │"
echo " │    datawall              daily report                    │"
echo " │    datawall start        start daemon                    │"
echo " │    datawall stop         stop daemon                     │"
echo " │    datawall restart      restart daemon                  │"
echo " │    datawall status       daemon status                   │"
echo " │    datawall reset        reset today's data              │"
echo " │    datawall watch        real-time monitor (htop-like)   │"
echo " │    datawall days N       aggregated report for N days    │"
echo " │    datawall limit        show/set daily limit            │"
echo " │    datawall connections  active TCP connections          │"
echo " └──────────────────────────────────────────────────────────┘"
